# docker/conf/php/supervisord.conf
; ─── supervisord 主进程 ──────────────────────────────────────────
[supervisord]
nodaemon=true
user=root
loglevel=warn
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
pidfile=/var/run/supervisord.pid

; ─── 管理接口（支持 supervisorctl 运维）────────────────────────
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ─── Octane（RoadRunner）守护进程 ────────────────────────────────
[program:octane]
command=php -d variables_order=EGPCS /srv/artisan octane:start --workers=4 --task-workers=6 --server=roadrunner --host=0.0.0.0 --rpc-port=7000 --port=8000 --max-requests=300
user=root
autostart=true
autorestart=true
startretries=3
; TERM 信号让 RoadRunner 优雅排干请求后退出
stopsignal=TERM
; 给足时间处理在途请求，超时后强杀
stopwaitsecs=30
; 确保 RoadRunner 子进程一并退出，避免僵尸进程
stopasgroup=true
killasgroup=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/octane.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
stdout_logfile=/var/log/supervisor/octane.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
; 最后启动、优先关闭
priority=999

; ─── Queue Worker ────────────────────────────────────────────────
[program:queue]
command=php /srv/artisan queue:work --max-jobs=1024 --timeout=60 --sleep=3 --tries=3
user=root
; 多进程并发消费，可按需调整 numprocs
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startretries=3
; QUIT 信号让 Laravel Queue 处理完当前 job 再退出
stopsignal=QUIT
; 必须 > --timeout(60)，给当前 job 留足完成时间
stopwaitsecs=65
stopasgroup=true
killasgroup=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/queue.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/supervisor/queue.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=997

; ─── 定时任务调度 ─────────────────────────────────────────────────
[program:schedule]
command=php /srv/artisan schedule:work
user=root
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/schedule.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=/var/log/supervisor/schedule.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
priority=996
